# photo-gallery.service
[Unit]
Description=Photo Gallery - Next.js Application
Documentation=https://github.com/frostiks777/1pgWebGallery
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/apps/photo-gallery

# Окружение
Environment="NODE_ENV=production"
Environment="PORT=3000"
Environment="NODE_OPTIONS=--max-old-space-size=4096"
EnvironmentFile=/var/www/apps/photo-gallery/.env.local

# Запуск через Bun (bunx next start)
ExecStart=/usr/local/bin/bunx next start -p 3000

# Перезапуск при падении
Restart=on-failure
RestartSec=10

# Таймауты
TimeoutStartSec=60
TimeoutStopSec=30

# Логирование
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=photo-gallery

# Ограничения ресурсов
LimitNOFILE=65535

# Безопасность
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
